# minios/ramfs/etc/init/minios.conf
# Copyright 2021 The ChromiumOS Authors
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

description     "Startup for MiniOS script"
author          "chromeos-core-services@google.com"

start on started dbus and started frecon and started shill
expect fork
respawn

pre-start script
  HAS_LIDSWITCH=0
  for f in /sys/class/input/*; do
    if [ "$(cat "${f}/name")" == "Lid Switch" ]; then
      HAS_LIDSWITCH=1
      cd "${f}/device"
      break
    fi
  done

  if [ ${HAS_LIDSWITCH} == "1" ]; then
    PATH=$(pwd -P)
    DEV=${PATH##*/}
    echo ${DEV} > ${PATH}/driver/unbind
  fi
end script

oom score never

exec /usr/bin/minios
